# mirror-android-phone-display-on-plugin/mirror-android-phone-display.service
[Unit]
Description=Mirror Android display on phone

[Service]
ExecStart=/usr/local/bin/mirror-android-phone-display.sh
Restart=on-failure